# sources/vendor/pkgr/scaling/systemd/gogs-web-PROCESS_NUM.service
[Unit]
StopWhenUnneeded=true
Requires=gogs-web.service
After=gogs-web.service

[Service]
Environment=PORT=PORT_NUM
ExecStart=/usr/bin/gogs run web
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

[Install]
WantedBy=gogs-web.service
